# src/streamdock_n3/_data/streamdock-n3.service
[Unit]
Description=Stream Dock N3 controller
After=graphical-session.target pipewire.service

[Service]
Type=simple
ExecStart=@BIN@/streamdock-n3
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
